# deploy/gentoo-zh-verify-bot.service
[Unit]
Description=Gentoo zh Telegram join-verification bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/gentoo-zh-verify-bot --config /etc/gentoo-zh-verify-bot/config.json
EnvironmentFile=/etc/gentoo-zh-verify-bot/bot.env
# persistent state dir (/var/lib/gentoo-zh-verify-bot) so in-progress verifications survive restarts
StateDirectory=gentoo-zh-verify-bot
Restart=on-failure
RestartSec=5

# resource cap: the bot is light (~50 MiB steady, brief peak while building the /pkg
# cache); cap it so a runaway can't pressure anything else on the host.
MemoryMax=512M

# hardening — the bot only needs outbound HTTPS + its StateDirectory, nothing else.
DynamicUser=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
ProtectProc=invisible
ProcSubset=pid
RestrictSUIDSGID=yes
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
LockPersonality=yes
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
UMask=0077

[Install]
WantedBy=multi-user.target
